# mail/.config/systemd/user/mbsync.service
[Unit]
Description=A job that syncs all mailboxes

[Service]
Type=simple
ExecStart=bash -c "PATH=$PATH:%h/.local/bin/ mbsync -a"
